# scripts/systemd/skyweaver-worker.service
[Unit]
Description=Sky Weaver Hub processing worker
After=skyweaver-api.service
PartOf=skyweaver.target

[Service]
Type=simple
User=skyweaver
Group=skyweaver
WorkingDirectory=/opt/skyweaver
EnvironmentFile=-/etc/skyweaver/skyweaver.env
ExecStart=/opt/skyweaver/backend/.venv/bin/python -m skyweaver.worker
Restart=on-failure
RestartSec=5

[Install]
WantedBy=skyweaver.target
