# scripts/systemd/skyweaver-api.service
[Unit]
Description=Sky Weaver Hub API
After=network-online.target
Wants=network-online.target
PartOf=skyweaver.target

[Service]
Type=simple
User=skyweaver
Group=skyweaver
SupplementaryGroups=video render input gpio i2c spi
WorkingDirectory=/opt/skyweaver
EnvironmentFile=-/etc/skyweaver/skyweaver.env
Environment=PYTHONPATH=/opt/skyweaver/backend
ExecStart=/opt/skyweaver/backend/.venv/bin/uvicorn skyweaver.main:app --host ${SKYWEAVER_HOST} --port ${SKYWEAVER_PORT}
Restart=on-failure
RestartSec=5

[Install]
WantedBy=skyweaver.target
