# scripts/systemd/skyweaver-capture.service
[Unit]
Description=Sky Weaver Hub capture daemon
After=skyweaver-api.service
PartOf=skyweaver.target

[Service]
Type=simple
User=skyweaver
Group=skyweaver
WorkingDirectory=/opt/skyweaver
EnvironmentFile=-/etc/skyweaver/skyweaver.env
Environment=PYTHONPATH=/opt/skyweaver/backend
ExecStart=/opt/skyweaver/backend/.venv/bin/python -m skyweaver.capture_daemon
Restart=on-failure
RestartSec=5

[Install]
WantedBy=skyweaver.target
